# inittab
::respawn:/usr/local/bin/start_syslogd
::respawn:/usr/local/bin/run_hsd
