# inittab
::respawn:/sbin/syslogd -n -O -
::respawn:/usr/local/bin/run_hsd
